# .config/systemd/user/nwbat.timer
[Unit]
Description=Battery Indicator Timer

[Timer]
OnUnitActiveSec=100
AccuracySec=10

[Install]
WantedBy=timers.target
