# .config/systemd/user/nwbat.timer
[Unit]
Description=Battery Indicator Timer

[Timer]
OnUnitActiveSec=10
AccuracySec=1

[Install]
WantedBy=timers.target
