# .config/systemd/user/nwbat.timer
[Unit]
Description=Battery Charging State Timer

[Timer]
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
