# .config/systemd/user/resume.target
[Unit]
Description=On Resume
After=suspend.target
StopWhenUnneeded=yes

[Install]
WantedBy=suspend.target
