# .config/systemd/user/sway.target
[Unit]
Description=Sway Window Manager
Documentation=man:systemd.special(7)
BindsTo=graphical-session.target
Requires=sway.service
