# .config/systemd/user/weather-update.timer
[Unit]
Description=Weather Timer

[Timer]
OnUnitActiveSec=1800

[Install]
WantedBy=timers.target
