# .config/systemd/user/weather.timer
[Unit]
Description=Weather Timer

[Timer]
OnBootSec=0
OnUnitActiveSec=1800

[Install]
WantedBy=timers.target
