# .config/systemd/user/weather.service
[Unit]
Description=Get Weather

[Service]
Type=oneshot
ExecStart=/home/hugo/.local/bin/weather Toulouse

[Install]
WantedBy=sway-session.target
